# templates/firewall/kvmfirewall.erb
#!/bin/bash
#
# init the firewall

### BEGIN INIT INFO
# Provides:          kvmfirewall
# Required-Start:    $network $syslog
# Required-Stop:     $network $syslog
# Should-Start:      $local_fs
# Should-Stop:       $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: ip tables firewall
# Description:       secure the host
### END INIT INFO
 
. /etc/environment

case "$1" in
	start)
		for SCRIPT in `ls /etc/firewall/*.sh`; do
			 $SCRIPT start
		done;
		;;
	stop)
		for SCRIPT in `ls /etc/firewall/*init-vlan*.sh`; do
			 $SCRIPT stop
		done;		
		/etc/firewall/000-init.sh stop
		;;
	restart)
		$0 stop
		$0 start
		;;
esac